# config/echo-daemon.service
[Unit]
Description=Echo Speech-to-Text Daemon
After=network.target sound.target
Wants=network.target

[Service]
Type=forking
User=%i
Group=%i
ExecStart=/home/nigel/echo/daemon/echo_daemon.py start
ExecStop=/home/nigel/echo/daemon/echo_daemon.py stop
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/tmp/echo_daemon.pid
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/tmp /home/%i/.config/echo /home/%i/.local/share/echo
PrivateTmp=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment
Environment=DISPLAY=:0
Environment=PULSE_SERVER=unix:/run/user/%i/pulse/native

[Install]
WantedBy=default.target
