# config/echo-daemon.service
[Unit]
Description=Echo Speech-to-Text Daemon
After=network.target sound.target
Wants=network.target

[Service]
Type=forking
# User and Group not needed for user services
ExecStart=/home/nigel/echo/daemon/echo_daemon.py start
ExecStop=/home/nigel/echo/daemon/echo_daemon.py stop
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/tmp/echo_daemon.pid
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Security settings (simplified for user service)
NoNewPrivileges=true

# Environment
Environment=DISPLAY=:0
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native

[Install]
WantedBy=default.target
